# layout/DEBIAN/postinst
#!/bin/sh
/bin/launchctl unload  /System/Library/LaunchDaemons/com.qwertyoruiop.uisettingsd || :
/bin/launchctl load  /System/Library/LaunchDaemons/com.qwertyoruiop.uisettingsd
:
